# process_monitor.service
[Unit]
Description=Monitor the process

[Service]
ExecStart=/usr/local/bin/process_monitor.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target 
